# data/kiran-sse-daemon.service.in
[Unit]
Description=Kiran SSE Daemon Service

[Service]
Type=dbus
BusName=com.kylinsec.Kiran.SSE
ExecStart=@SSE_INSTALL_BINDIR@/kiran-sse-daemon

[Install]
WantedBy=multi-user.target
