# data/kiran-ssr-daemon.service.in
[Unit]
Description=Kiran SSR Daemon Service

[Service]
Type=dbus
BusName=com.kylinsec.Kiran.SSR
ExecStart=@SSR_INSTALL_BINDIR@/kiran-ssr-daemon

[Install]
WantedBy=multi-user.target
