# data/ks-sc-daemon.service.in
[Unit]
Description=KylinSec SC Daemon Service

[Service]
Type=dbus
BusName=com.kylinsec.SC
ExecStart=@SC_INSTALL_BINDIR@/ks-sc-daemon

[Install]
WantedBy=multi-user.target
