# data/ks-ssr-daemon.service.in
[Unit]
Description=KylinSec SC Daemon Service
Requires=ksl-daemon.service
After=ksl-daemon.service

[Service]
Type=dbus
BusName=com.kylinsec.SSR
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ks-ssr-daemon

[Install]
WantedBy=multi-user.target
