# data/ks-ssr-daemon.service.in
[Unit]
Description=KylinSec SSR Daemon Service

[Service]
Type=dbus
BusName=com.kylinsec.SSR
ExecStart=@SSR_INSTALL_BINDIR@/ks-ssr-daemon

[Install]
WantedBy=multi-user.target
